# androidFarm/supervisord.conf
[supervisord]
nodaemon=true

[program:adb]
command=adb start-server

[program:jenkinsSlave]
command=su androidfarm -c "/bin/bash /home/androidfarm/files/androidFarm.sh"
